# default/systemd/wayland/display-manager.path
[Unit]
Description=Wait for Xorg socket
Requires=display-manager-run.service
After=display-manager-run.service

[Path]
PathExists=/usr/bin/Xorg
