# default/wayland/display-manager.path
[Unit]
Description=Wait for display daemon socket

[Path]
PathExists=/tmp/wayland-0
