# default/wayland/display-manager.path
[Unit]
Description=Wait for display daemon socket
Requires=display-manager-run.service
After=display-manager-run.service

[Path]
PathExists=/tmp/wayland-0
