# default/wayland/display-manager-run.service
[Unit]
Description=Display manager run service
Requires=enlightenment.service
After=enlightenment.service

[Service]
Type=oneshot

[Install]
WantedBy=graphical.target
